# ⁄etc⁄systemd⁄system⁄spegnimento-negozio.service
[Unit]
Description=Chiusura applicativi e spegnimento di sistema

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl isolate poweroff.target
