# ⁄etc⁄systemd⁄system⁄spegnimento-negozio.timer
[Unit]
Description=Timer per spegnimento negozio alle 21:00

[Timer]
OnCalendar=*-*-* 21:00:00
Persistent=true

[Install]
WantedBy=timers.target
